# init.d.in
#!/bin/bash
#
# host-secrets
# This shell script takes care of starting and stopping taskcluster-host-secrets
#
# chkconfig: - 58 74
# description: host-secrets is the taskcluster-host-secrets daemon.

# Source function library.
. /etc/init.d/functions

prog=host-secrets
lockfile=/var/lock/subsys/$prog
conffile=/etc/$prog.conf

[ -f $conffile ] || exit 6
. $conffile
pid=$(fuser -n tcp $PORT)

start() {
  [ "$EUID" != "0" ] && exit 4
  [ -f <<<BINDIR>>>/../../bin/start.sh ] || exit 5

  if [ "$pid" != "" ] && ps -p $pid > /dev/null; then
    echo -n $"$prog is already running."
    return 0
  else
    echo -n $"Starting $prog: "
    daemon $prog $OPTIONS
    exitcode=$?
    echo
    [ $exitcode -eq 0 ] && touch $lockfile
    return $exitcode
  fi
}

stop() {
  [ "$EUID" != "0" ] && exit 4
  if [ "$pid" != "" ] && ps -p $pid > /dev/null; then
    echo -n $"Shutting down $prog: "
    kill $(fuser -n tcp $PORT)
    exitcode=$?
    echo
    [ $exitcode -eq 0 ] && rm -f $lockfile
    return $exitcode
  else
    echo -n $"$prog is not running."
    return 0
  fi
}

get_status() {
  if [ "$pid" != "" ] && ps -p $pid > /dev/null; then
    return 0
  elif [ -f $lockfile ]; then
    return 2
  else
    return 3
  fi
}

case "$1" in
  start)
  start
  ;;
  stop)
  stop
  ;;
  status)
  exit $(get_status)
  ;;
  restart)
  stop
  start
  ;;
  *)
  echo $"Usage: $0 {start|stop|status|restart}"
  exit 2
esac
